# bin/systemd/ndm.target
[Unit]
Description=Start NDM Application
StartLimitIntervalSec=1s
StartLimitBurst=1000
After=network.target
[Install]
WantedBy=multi-user.target
